# infra/backup/systemd/questionshub-restic-check.timer
[Unit]
Description=Weekly restic check

[Timer]
OnCalendar=Sun *-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
